# files/leodis.conf
start on runlevel [2345]
stop on starting rc RUNLEVEL=[16]

exec nodejs /var/node/leodis/app.js
